# btc_usdt_bot_usb_copy/bot_manager.service
[Unit]
Description=BTC/USDT Trading Bot Dashboard Service
After=network.target

[Service]
User=root
Group=root
# The setup script will replace THIS_PATH with the actual directory
WorkingDirectory=THIS_PATH
ExecStart=/bin/bash THIS_PATH/start_dashboard.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
